# detections/004_scheduled_task_creation_persistence.yaml
id: 004-scheduled-task-persistence
name:  Scheduled Task Creation (Persistence)
description: Detects creation of new scheduled tasks which are commonly used for persistence and execution.
severity: Medium
status: Available
requiredDataConnectors:
  - connectorId: AzureMonitorWindowsAgents
    dataTypes: [SecurityEvent]
tactics: [Persistence, Execution]
techniques: [T1053.005]
queryFrequency: 1h
queryPeriod: 1h
triggerOperator: gt
triggerThreshold: 0
kind: Scheduled
query: |-
  SysmonParsed
  | where EventID == 1
  | where Image endswith "schtasks.exe"
  | where CommandLine has "/create"
  | project TimeGenerated, Computer, User, Image, CommandLine, ParentImage
version: 1.0.0
tags:
  - T1053.005
  - Persistence
